# systemd/made-by-ai-autoloop.service.template
[Unit]
Description=Made by AI OpenCode Auto Loop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_ROOT__
Environment=HOME=%h
Environment=PATH=/usr/local/bin:/usr/bin:/bin:%h/.local/bin
EnvironmentFile=__PROJECT_ROOT__/.env.local
ExecStart=/bin/bash __PROJECT_ROOT__/scripts/autoloop.sh
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
